# src/pygenkit/templates/python-cli/debian/{{module_name}}.service.jinja
[Unit]
Description={{ project.name }}
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/{{ module_name }}
Restart=on-failure
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
